# hardware/coffee-light.service
[Unit]
Description=Websocket listener for coffee.waeco-soft.com
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/coffee-light/mySocket.py
WorkingDirectory=/opt/coffee-light
Restart=always
User=root

[Install]
WantedBy=multi-user.target
